# health-monitor.service
[Unit]
Description=GoStorm Health Monitor Dashboard
After=network.target gostream.service

[Service]
Type=simple
# Change User to your system account
User=YOUR_USER
# Change to your GoStream install directory
WorkingDirectory=/home/pi/GoStream
ExecStart=/usr/bin/python3 scripts/health-monitor.py
Restart=always
RestartSec=5
TimeoutStopSec=5
Environment="PYTHONUNBUFFERED=1"

StandardOutput=append:logs/health-monitor.log
StandardError=append:logs/health-monitor.log

[Install]
WantedBy=multi-user.target